# bin/build/tesseract-ocr.conf
description "A python OCR utility for scanned PDF or image document files.  This upstart program will handle the daemonization process for the python script."
author "Corey Farmer <corey.farmer@emerus.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec python /home/cfarmer/Coding/com.sadmicrowave.tesseract/tesseract-daemon.py -v